# vendor/cookbooks/osx-bootstrap/templates/default/ssh-homebrew.openssh.ssh-agent.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- This file is managed by Chef. Your changes will be overwritten. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>homebrew.openssh.ssh-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string><%= (prefix + "bin/ssh-agent").to_s.encode(xml: :text) -%></string>
        <string>-l</string>
    </array>
    <key>Sockets</key>
    <dict>
        <key>Listeners</key>
        <dict>
            <key>SecureSocketWithKey</key>
            <string>SSH_AUTH_SOCK</string>
        </dict>
    </dict>
    <key>EnableTransactions</key>
    <true/>
</dict>
</plist>
